# backend/assortiShop.service
[Unit]
Description=AssortiShop Backend and Database
After=docker.service
Requires=docker.service
Documentation=file:///root/sofar5/SELECTEL_DEPLOYMENT.md

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/root/sofar5/backend
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
